# systemd/root/kill-sshfs-on-suspend.service
# https://wiki.archlinux.org/index.php/Power_management

[Unit]
Before=sleep.target

[Service]
Type=simple
ExecStart=-/usr/bin/pkill sshfs

[Install]
WantedBy=sleep.target
