# services/user/nm-applet.service
[Unit]
Description=nm-applet

[Service]
ExecStart=nm-applet

[Install]
WantedBy=user-applications.target
