# services/user/pasystray.service
[Unit]
Description=pasystray

[Service]
ExecStart=pasystray

[Install]
WantedBy=user-applications.target
