# systemd/user/blueman-applet.service
[Unit]
Description=blueman-applet

[Service]
ExecStart=blueman-applet

[Install]
WantedBy=user-applications.target
