# systemd/user/copyq.service
[Unit]
Description=copyq

[Service]
ExecStart=copyq

[Install]
WantedBy=user-applications.target
