# systemd/user/discord.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=discord
ExecStop=kill $MAINPID

[Install]
WantedBy=graphical-session.target
