# systemd/user/discord.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=discord

[Install]
WantedBy=graphical-session.target
