# systemd/user/earlyoom.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=earlyoom

[Install]
WantedBy=graphical-session.target
