# systemd/user/evscript.service
[Unit]
Description=evscript

[Service]
ExecStart=bash -c "fd '^.*kbd$' /dev/input/by-path | xargs sudo evscript -f /home/cjbassi/config/other/evscript.dyon -d"
ExecStop=sudo pkill evscript

[Install]
WantedBy=user-applications.target
