# systemd/user/kdeconnect-indicator.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=kdeconnect-indicator

[Install]
WantedBy=graphical-session.target
