# systemd/user/mako.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=mako --default-timeout 7500

[Install]
WantedBy=graphical-session.target
