# systemd/user/mako.service
[Unit]
Description=mako

[Service]
ExecStart=mako

[Install]
WantedBy=user-applications.target
