# systemd/user/mako.service
[Unit]
Description=mako

[Service]
ExecStart=mako --default-timeout 7500

[Install]
WantedBy=user-applications.target
