# systemd/user/mons.service
[Unit]
Description=mons

[Service]
ExecStart=mons -a -x /home/cjbassi/config/other/mons_daemon_script

[Install]
WantedBy=user-applications.target
