# systemd/user/nm-applet.service
[Unit]
Description=nm-applet

[Service]
ExecStart=nm-applet --indicator

[Install]
WantedBy=user-applications.target
