# systemd/user/recover-youtube-videos.timer
[Unit]
Description=Run recover-youtube-videos daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
