# systemd/user/redshift-gtk.service
[Unit]
Description=redshift-gtk

[Service]
ExecStart=redshift-gtk -m wayland

[Install]
WantedBy=user-applications.target
