# systemd/user/riot-desktop.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=riot-desktop

[Install]
WantedBy=graphical-session.target
