# systemd/user/udiskie.service
[Unit]
Description=udiskie

[Service]
ExecStart=udiskie --no-automount --tray

[Install]
WantedBy=user-applications.target
