# systemd/user/unclutter.service
[Unit]
Description=unclutter

[Service]
ExecStart=unclutter --timeout 1

[Install]
WantedBy=user-applications.target
