# systemd/user/variety.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=variety

[Install]
WantedBy=graphical-session.target
