# systemd/user/variety.service
[Unit]
Description=variety

[Service]
ExecStart=variety

[Install]
WantedBy=user-applications.target
