# systemd/user/waybar.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=waybar

[Install]
WantedBy=graphical-session.target
