# systemd/user/waybar.service
[Unit]
Description=waybar

[Service]
ExecStart=waybar

[Install]
WantedBy=user-applications.target
