# systemd/user/xcape.service
[Unit]
Description=xcape

[Service]
ExecStart=xcape -d

[Install]
WantedBy=user-applications.target
